# supervisord.conf
[supervisord]
pidfile=/wayback/supervisord.pid
loglevel=warn
logfile=/dev/stdout
logfile_maxbytes=0
nodaemon=true
nocleanup=false
silent=true

[program:dbus]
;command=/usr/bin/dbus-daemon --system --nofork
command=/usr/bin/dbus-daemon --session
priority=1
numprocs=1
autostart=true
autorestart=true
startsecs=1
stderr_logfile=/dev/stderr
stdout_logfile=/dev/stdout
stderr_logfile_maxbytes=0
stdout_logfile_maxbytes=0

[program:xvfb]
command=Xvfb -ac :99 -screen 0 1280x1024x16
priority=10
numprocs=1
autostart=true
autorestart=true
startsecs=10
startretries=2147483647
stderr_logfile=/dev/stderr
stdout_logfile=/dev/stdout
stderr_logfile_maxbytes=0
stdout_logfile_maxbytes=0

[program:starter]
command=/bin/sh -c 'sleep 1 && set -o pipefail && /bin/sh /starter.sh'
priority=20
numprocs=1
autostart=true
autorestart=false
startsecs=1
startretries=1
stderr_logfile=/dev/stderr
stdout_logfile=/dev/stdout
stderr_logfile_maxbytes=0
stdout_logfile_maxbytes=0

;[program:chromium]
;command=/bin/sh -c 'sleep 10 && set -o pipefail && /usr/bin/chromium-browser --disable-gpu --no-sandbox --remote-debugging-address=0.0.0.0 --remote-debugging-port=9222 --disable-dev-shm-usage --disable-translate --disable-extensions --disable-software-rasterizer --ignore-certificate-errors --disable-sync --use-fake-ui-for-media-stream --use-gl=swiftshader --hide-scrollbars --disable-default-apps --mute-audio --no-default-browser-check --no-first-run --disable-renderer-backgrounding --disable-backgrounding-occluded-windows --disable-background-timer-throttling --force-fieldtrials="*BackgroundTracing/default/" --js-flags="--max_old_space_size=%(ENV_GC_MEM)s" --user-data-dir="%(ENV_CHROMIUM_WORKSPACE)s"'
;priority=30
;numprocs=1
;autostart=true
;autorestart=true
;startsecs=1
;stderr_logfile=/dev/stderr
;stdout_logfile=/dev/stdout
;stderr_logfile_maxbytes=0
;stdout_logfile_maxbytes=0

[program:tor]
command=/usr/bin/tor -f /etc/tor/torrc --quiet --SocksPort 9050
priority=40
numprocs=1
autostart=true
autorestart=true
startsecs=1
stderr_logfile=/dev/stderr
stdout_logfile=/dev/stdout
stderr_logfile_maxbytes=0
stdout_logfile_maxbytes=0

[program:wayback]
command=/bin/sh -c 'set -o pipefail && /bin/sh /entrypoint.sh'
priority=50
numprocs=1
autostart=true
autorestart=true
startsecs=10
startretries=2147483647
stderr_logfile=/dev/stderr
stdout_logfile=/dev/stdout
stderr_logfile_maxbytes=0
stdout_logfile_maxbytes=0

[program:forward]
command=/usr/bin/socat %(ENV_SOCAT_OPTIONS)s tcp-listen:%(ENV_PORT)s,reuseaddr,fork,nonblock,nodelay tcp:localhost:%(ENV_WAYBACK_TOR_LOCAL_PORT)s
priority=60
numprocs=1
autostart=true
autorestart=true
startsecs=1
stderr_logfile=/dev/stderr
stdout_logfile=/dev/stdout
stderr_logfile_maxbytes=0
stdout_logfile_maxbytes=0

[program:cleaner]
command=/bin/sh -c 'sleep 60 && set -o pipefail && /bin/sh /cleaner.sh'
numprocs=1
autostart=true
autorestart=true
startsecs=10
startretries=2147483647
stderr_logfile=/dev/stderr
stdout_logfile=/dev/stdout
stderr_logfile_maxbytes=0
stdout_logfile_maxbytes=0
